# v4_execution.service
[Unit]
Description=BitMEX V4 Cascade Dip-Buy Service (TESTNET)
After=network.target

[Service]
User=meidie
WorkingDirectory=/home/meidie/BitMexBot
ExecStart=/home/meidie/BitMexBot/venv/bin/python -u v4_execution.py
Restart=always
RestartSec=30
StandardOutput=append:/home/meidie/BitMexBot/logs/v4.log
StandardError=append:/home/meidie/BitMexBot/logs/v4.log
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
